# voice-cli.service
[Unit]
Description=Voice CLI Daemon
After=network.target

[Service]
Type=simple
User=snehit
WorkingDirectory=/home/snehit/projects/voice-cli
ExecStart=/home/snehit/projects/voice-cli/node_modules/.bin/bun run /home/snehit/projects/voice-cli/index.ts start --no-supervisor
Restart=always
RestartSec=5
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/home/snehit/.amp/bin:/home/snehit/.local/zed.app/bin:/home/snehit/.bun/bin:/usr/bin:/home/snehit/.atuin/bin:/home/linuxbrew/.linuxbrew/bin:/home/linuxbrew/.linuxbrew/sbin:/home/snehit/go/bin:/home/snehit/.local/bin:/home/snehit/.npm-global/bin:/usr/lib64/ccache:/home/snehit/.cargo/bin:/usr/local/bin:/usr/bin:/usr/local/sbin:/usr/sbin:/var/lib/snapd/snap/bin
Environment=DISPLAY=:0
Environment=XAUTHORITY=
Environment=WAYLAND_DISPLAY=wayland-1

[Install]
WantedBy=default.target
